# shard-enterprise.service
[Unit]
Description=SHARD Enterprise SIEM
After=network.target remote-fs.target

[Service]
Type=simple
User=user
WorkingDirectory=/mnt/c/Users/user/PycharmProjects/Shard
Environment="PYTHONPATH=/mnt/c/Users/user/PycharmProjects/Shard"
Environment="HOME=/home/user"
Environment="MPLCONFIGDIR=/tmp/matplotlib"
EnvironmentFile=/mnt/c/Users/user/PycharmProjects/Shard/.env
ExecStartPre=/bin/bash -c 'until [ -f /mnt/c/Users/user/PycharmProjects/Shard/run_shard.py ]; do sleep 2; done'
ExecStartPre=/bin/sleep 5
ExecStart=/home/user/.venv-linux/bin/python3 /mnt/c/Users/user/PycharmProjects/Shard/run_shard.py --no-capture
Restart=on-failure
RestartSec=30
StartLimitBurst=10
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
